# media-diary.service
[Unit]
Description=Media Diary FastAPI app
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/media-diary
EnvironmentFile=/home/pi/media-diary/.env
Environment=REPO_ROOT=/home/pi/media-diary
ExecStart=/home/pi/media-diary/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8765
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
